# antigravity-dashboard.service
[Unit]
Description=Antigravity Dashboard
After=network.target

[Service]
Type=simple
User=mateusz
WorkingDirectory=/home/mateusz/MCPTA
Environment="PATH=/home/mateusz/MCPTA/venv/bin:/usr/bin:/bin"
EnvironmentFile=/home/mateusz/MCPTA/.env
ExecStart=/home/mateusz/MCPTA/venv/bin/python /home/mateusz/MCPTA/dashboard.py
StandardOutput=append:/tmp/dashboard.log
StandardError=append:/tmp/dashboard.log
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
